# resources/pulse.timer
[Unit]
Description=Pulse daemon

[Timer]
OnCalendar=*:0/5
 
[Install]
WantedBy=timers.target
