# resources/pvoutput.timer
[Unit]
Description=Pulsemeter upload to PVoutput

[Timer]
OnCalendar=*:1/5
 
[Install]
WantedBy=timers.target
