# resources/pvoutput.service
[Unit]
Description=Pulsemeter upload to PVOutput

[Service]
ExecStart=/scratch/git/pulse/resources/pvoutput.sh
User=rrdcached
Group=rrdcached

[Install]
WantedBy=multi-user.target
